# generators/bluepill_init.conf.erb
# bluepill - process monitor - 
# Rails application: <%= "#{application}" %>
# Environment: <%= "#{rails_env}" %>
#
# simple process monitoring tool

description     "Simple process monitoring tool for <%= "#{application}" %> - <%= "#{rails_env}" %>"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

pre-start script
  mkdir -p /var/run/bluepill
  chown -R deploy:www-data /var/run/bluepill
end script

#env PATH=/opt/ruby-enterprise-1.8.7-2010.01/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games
script
  chdir <%= "#{current_path}" %>
  exec su deploy -c 'bundle exec bluepill load <%= "#{bluepill_remote_config}" %>'
end script
