# generators/bluepill_init.conf.erb
<% ruby_version     = "1.9.3" if ruby_version.nil? %>
# bluepill - process monitor - 
# Rails application: <%= "#{application}" %>
# Environment: <%= "#{rails_env}" %>
#
# simple process monitoring tool

description     "Simple process monitoring tool for <%= "#{application}" %> - <%= "#{rails_env}" %>"

#start on runlevel [2345]
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

expect daemon
# respawn

pre-start script
  echo "Pre-start <%= "#{application}" %>"
  mkdir -p /var/run/bluepill
  chown -R <%= "#{user}:#{group}" %> /var/run/bluepill
end script

pre-stop script
  echo "Pre-stop <%= "#{application}" %>"
  exec su -c "source '/home/leif/.rvm/scripts/rvm';rvm <%= "#{ruby_version}" %>; cd <%= "#{current_path}" %>; RAILS_ENV=<%= "#{rails_env}" %> bundle exec bluepill normedia-download_staging stop --no-privileged" <%= "#{user}" %>
end script

script
  echo "Starting <%= "#{application}" %>"
  exec su -c "source '/home/<%= "#{user}" %>/.rvm/scripts/rvm';rvm <%= "#{ruby_version}" %>; cd <%= "#{current_path}" %>; RAILS_ENV=<%= "#{rails_env}" %> bundle exec bluepill load <%= "#{bluepill_remote_config}" %> --no-privileged" <%= "#{user}" %>
end script


